# backend/ops/macos/com.yasargold.backup-postgres.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yasargold.backup-postgres</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd "$HOME/yasargold" && ./backend/backup_postgres.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/yasargold-backup-postgres.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/yasargold-backup-postgres.err.log</string>
</dict>
</plist>
